# dvd2m4v.service
[Unit]
Description=Encode DVDs to m4vs.
Wants=dvd2m4v.timer

[Service]
ExecStart=dvd2m4v
TimeoutStartSec=3600